# supervisord.conf
[supervisord]
nodaemon=true

[program:streamlit]
command=streamlit run /app/src/main.py
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:fastapi]
command=uvicorn src.api:app --host 0.0.0.0 --port 8000 --loop uvloop --http httptools
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr